# rancid/rancid.cron
0/30 * * * * rancid /usr/lib/rancid/bin/rancid-run
